# deploy/sentient-analysis-queue.service
[Unit]
Description=Sentient Trader Analysis Queue Processor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/sentient-trader
Environment="PATH=/root/sentient-trader/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/root/sentient-trader"
EnvironmentFile=/root/sentient-trader/.env
ExecStart=/bin/bash -c '/root/sentient-trader/venv/bin/python3 -u windows_services/runners/run_analysis_queue_processor.py 2>&1 | tee -a /root/sentient-trader/logs/analysis_queue_service.log'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
